# packages/shared-state-nodes_and_links/files/etc/uci-defaults/shared-state-bat_nodes_and_links
#!/bin/sh

. /lib/libremesh.sh

unique_append \
	'*/5 * * * * ((sleep $(($RANDOM % 120)); shared-state sync nodes_and_links &> /dev/null)&)'\
	/etc/crontabs/root
