# packages/first-boot-wizard/files/etc/init.d/firstbootwizard
#!/bin/sh /etc/rc.common

USE_PROCD=1
NAME=firstbootwizard
PROG=/bin/firstbootwizard

start_service() {
	procd_open_instance
	procd_set_param command "$PROG"
	procd_close_instance
}

rescan() {
	service_reload $PROG
}