# packages/tests-scheduler/files/etc/uci-defaults/51_tests-scheduler-at
#!/bin/sh

hour=$(awk 'BEGIN{srand();print int(rand()*24)}')
minute=$(awk 'BEGIN{srand();print int(rand()*60)}')

unique_append()
{
	grep -qF "$1" "$2" || echo "$minute $hour $1" >> "$2"
}

unique_append \
	'* * * ((shared-state sync tests-scheduler-night &> /dev/null; /bin/tests-scheduler-at; shared-state sync tests-scheduler-night &> /dev/null)&)'\
	/etc/crontabs/root

exit 0
