# packages/lime-auto-ap-watchping/files/etc/uci-defaults/91_autoAP
#!/bin/sh
uci -q get lime.network.autoap_enabled || uci set lime.network.autoap_enabled=0
uci -q get lime.network.autoap_hosts || uci set lime.network.autoap_hosts="8.8.8.8 8.8.4.4"
exit 0
