# packages/shared-state-bat_hosts/files/etc/hotplug.d/iface/shared-state-bat_hosts.sh
#!/bin/sh
[ "x$ACTION" == "xifup" ] && ((sleep 30; shared-state-publish_bat_hosts; shared-state-async sync bat-hosts)&)
