# packages/shared-state-bat_hosts/files/etc/hotplug.d/iface/shared-state-bat_hosts.sh
#!/bin/sh
[ "x$ACTION" == "xifup" ] && ((/usr/bin/shared-state-publish_bat_hosts)&)
