# packages/shared-state-date/files/etc/uci-defaults/shared-state-date-cron
#!/bin/sh

uci set system.ntp.enabled=1
uci set system.ntp.enable_server=1
uci commit system

unique_append()
{
	grep -qF "$1" "$2" || echo "$1" >> "$2"
}

unique_append \
	'*/4 * * * * ((sleep $(($RANDOM % 90)); shared-state sync date)&)'\
	/etc/crontabs/root

unique_append \
	'*/2 * * * * ((sleep $(($RANDOM % 30)); /etc/shared-state/publishers/date)&)'\
	/etc/crontabs/root
