# packages/lime-webui/files/etc/profile.d/20-system_notes.sh
#!/bin/sh
echo ""
echo "=== System Notes ====================================="
[ -f /etc/banner.notes ] && cat /etc/banner.notes
echo "======================================================"
echo ""
