# packages/lime-webui/files/etc/profile.d/system_notes.sh
#!/bin/sh

echo System Notes
echo ============
cat /etc/banner.notes
echo
