# packages/shared-state-node_info/files/etc/uci-defaults/shared-state_node_info_cron
#!/bin/sh

unique_append()
{
	grep -qF "$1" "$2" || echo "$1" >> "$2"
}

uci set shared-state.node_info=dataType
uci set shared-state.node_info.name='node_info'
uci set shared-state.node_info.scope='community'
uci set shared-state.node_info.ttl='2400'
uci set shared-state.node_info.update_interval='33'
uci commit shared-state
