# packages/shared-state-bat_links_info/files/etc/uci-defaults/shared-state_bat_links_info_cron
#!/bin/sh

unique_append()
{
	grep -qF "$1" "$2" || echo "$1" >> "$2"
}
## information generator datatype
uci set shared-state.bat_links_info=dataType
uci set shared-state.bat_links_info.name='bat_links_info'
uci set shared-state.bat_links_info.scope='community'
uci set shared-state.bat_links_info.ttl='2400'
uci set shared-state.bat_links_info.update_interval='30'
uci commit shared-state

## Reference state data type
uci set shared-state.bat_links_info_ref=dataType
uci set shared-state.bat_links_info_ref.name='bat_links_info_ref'
uci set shared-state.bat_links_info_ref.scope='community'
uci set shared-state.bat_links_info_ref.ttl='2400'
uci set shared-state.bat_links_info_ref.update_interval='50'
