# packages/shared-state-wifi_links_info/files/etc/uci-defaults/shared-state_wifi_links_info
#!/bin/sh
unique_append()
{
	grep -qF "$1" "$2" || echo "$1" >> "$2"
}

unique_append \
	'*/5 * * * * ((sleep $(($RANDOM % 120)); shared-state sync wifi_links_info &> /dev/null)&)'\
	/etc/crontabs/root
