# packages/shared-state-ref_state_commons/files/etc/uci-defaults/90-shared-state-ref-persistence
#!/bin/sh
#safe upgrade looks for lime-autogen, wich is generated from defaults 
uci set lime-defaults.system.keep_on_upgrade="$(uci get lime-defaults.system.keep_on_upgrade) shared-state-ref-state"
uci commit lime-defaults
touch /tmp/done
exit 0
