# packages/shared-state-ref_state_commons/files/etc/uci-defaults/90-shared-state-ref-persistence
#!/bin/sh
uci set lime-autogen.system.keep_on_upgrade="$(uci get lime-autogen.system.keep_on_upgrade) shared-state-ref-state"
uci commit lime-autogen
exit 0
