# packages/batman-adv-auto-gw-mode/files/etc/uci-defaults/network_add_lan_dhcp
#!/bin/sh

uci batch <<EOF
set network.lan_dhcp=interface
set network.lan_dhcp.proto=dhcp
commit network
EOF
