# packages/ubus-lime-fbw/files/etc/init.d/ubus-lime-fbw
#!/bin/sh /etc/rc.common

START=99

USE_PROCD=1
NAME=fbw
PROG=/usr/libexec/daemon/lime-fbw

start_service() {
	procd_open_instance
	procd_set_param command "$PROG"
	procd_close_instance
}

stop() {
	service_stop $PROG
}

reload() {
	service_reload $PROG
}