# packages/dnsmasq-distributed-hosts/files/etc/uci-defaults/90_dnsmasq-distributed-hosts
#!/bin/sh

uci set alfred.alfred.disabled=
uci commit alfred

exit 0
