# packages/lime-webui/files/etc/uci-defaults/inject-luci-notify-lime-env.sh
#!/bin/bash

echo "<link rel='stylesheet' href='/luci-static/lime/notify-lime-env.css'>" >> /usr/lib/lua/luci/view/footer.htm
echo "<script src='/luci-static/lime/notify-lime-env.js'></script>" >> /usr/lib/lua/luci/view/footer.htm
