# packages/lime-system/files/etc/uci-defaults/90_lime-banner
#!/bin/sh

cat <<EOF >> /etc/banner
  ___   __ __                   _______             __
 |   |_|__|  |--.----.-----. _ |   |   |-----.-----|  |--.
 |     |  |  _  |   _|  -__||_||       |  -__|__ --|     |
 |_____|__|_____|__| |_____|   |__|_|__|_____|_____|__|__|

 ---------------------------------------------------------
 %LIME_DESCRIPTION%
 ---------------------------------------------------------
 http://libre-mesh.org
 ---------------------------------------------------------
EOF

exit 0
