# packages/lime-system/files/etc/uci-defaults/90_lime-config
#!/bin/sh
echo "Configuring LiMe for first boot..." | tee /tmp/lime-config.log
uci add_list uhttpd.main.listen_https="[::]:443"
uci add_list uhttpd.main.listen_http="[::]:80"
lua /usr/bin/lime-config | tee -a /tmp/lime-config.log
echo "###### uci changes generated by lime-config run: ######" >> /tmp/lime-config.log
uci changes >> /tmp/lime-config.log
uci commit
exit 0
