# packages/lime-system/files/etc/uci-defaults/91_lime-config
#!/bin/sh

mkdir -p /tmp/log/
echo "Configuring LiMe for first boot..." | tee /tmp/log/lime-config.log
lua /usr/bin/lime-config --no-commit | tee -a /tmp/log/lime-config.log
echo "###### uci changes generated by lime-config run: ######" >> /tmp/log/lime-config.log
uci changes >> /tmp/log/lime-config.log
uci commit
