# packages/lime-system/files/etc/uci-defaults/92_add-lime-repos
#!/bin/sh

. /etc/os-release
[ -f /etc/lime_release ] && . /etc/lime_release

feeds_file="/etc/opkg/limefeeds.conf"

[ -f "$feeds_file" ] && {
  echo "LibreMesh opkg feeds already defined - skipping"
  exit 0
}

[ -z "$OPENWRT_ARCH" ] && {
  echo "Release information not available, skipping opkg configuration"
  exit 0
}

[ "$LIME_BRANCH" == "master" ] || [ "$LIME_BRANCH" == "develop" ] || \
		[ "$LIME_BRANCH" == "staging" ] || [ -z "$LIME_RELEASE" ] && {
	base_url="http://snapshots.libremesh.org/packages/$OPENWRT_ARCH";
	key_name="7546f62c3d9f56b1"
	key_content="RWR1RvYsPZ9WsTsmoVajHX9dyl2wL7grjcfFua1q8A99RWr2gF94lRdJ"
} || {
	base_url="http://repo.libremesh.org/releases/$LIME_RELEASE/packages/$OPENWRT_ARCH"
	key_name="a71b3c8285abd28b"
	key_content="RWSnGzyChavSiyQ+vLk3x7F0NqcLa4kKyXCdriThMhO78ldHgxGljM/8"
}

echo "Configuring official LibreMesh opkg feeds"
echo "src/gz libremesh $base_url/libremesh" > /etc/opkg/limefeeds.conf
echo "untrusted comment: signed by libremesh.org key $key_name" > "/etc/opkg/keys/$key_name"
echo "$key_content" >> "/etc/opkg/keys/$key_name"
