# packages/lime-system/files/etc/uci-defaults/92_add-lime-repos
[ ! -f /etc/lime_release ] && {
	echo "LibreMesh version not found in /etc/lime_release - skipping" && exit 0
}

[ ! -f /bin/apk ] && [ ! -f /bin/opkg ] &&
	echo "No package manager found - skipping" && exit 0

. /etc/os-release
. /etc/lime_release

feed_host='http://feed.libremesh.org'
feed_ref=$([ $LIME_CODENAME == "development" ] &&
	echo "master" || echo "$LIME_RELEASE")
ow_branch=$( [ $VERSION == 'SNAPSHOT' ] &&
	echo "openwrt-main" || echo openwrt-${VERSION:0:5})

pkg_manager=apk
pkg_index=/packages.adb
feed_lime=/etc/apk/repositories.d/limefeeds.list
feed_keyname=libremesh.pem
feed_key="-----BEGIN PUBLIC KEY-----
MFkwEwYHKoZIzj0CAQYIKoZIzj0DAQcDQgAEdFJZ2qVti49Ol8LJZYuxgOCLowBS
8bI86a7zqhSbs5yon3JON7Yee7CQOgqwPOX5eMALGOu8iFGAqIRx5YjfYA==
-----END PUBLIC KEY-----"

if [ -f /bin/opkg ]; then
	pkg_manager=opkg
	pkg_index=''
	feed_lime=/etc/opkg/limefeeds.conf
	feed_keyname=a71b3c8285abd28b
	feed_key="untrusted comment: signed by libremesh.org key a71b3c8285abd28b
RWSnGzyChavSiyQ+vLk3x7F0NqcLa4kKyXCdriThMhO78ldHgxGljM/8"
	feed_prefix="src/gz "
fi

packages_url="$feed_host/$feed_ref/$ow_branch/x86_64$pkg_index";
arch_packages_url="$feed_host/$feed_ref/$ow_branch/${OPENWRT_ARCH}$pkg_index";
profiles_url="$feed_host/profiles/$ow_branch/x86_64$pkg_index";

[ $feed_prefix ] && {
	packages_url="${feed_prefix}libremesh $packages_url"
	arch_packages_url="${feed_prefix}libremesh_arch_packages $arch_packages_url"
	profiles_url="${feed_prefix}profiles $profiles_url"
}

echo "Configuring official LibreMesh $pkg_manager feeds"
cat << EOF > $feed_lime
$packages_url
$arch_packages_url
$profiles_url
EOF
echo "$feed_key" > "/etc/$pkg_manager/keys/$feed_keyname"
