# packages/lime-system/files/etc/uci-defaults/92_add-lime-repos
#!/bin/sh

. /etc/os-release
[ -f /etc/lime_release ] && . /etc/lime_release

feeds_file="/etc/opkg/limefeeds.conf"

[ -f "$feeds_file" ] && {
  echo "LibreMesh opkg feeds already defined - skipping"
  exit 0
}

ARCH="$OPENWRT_ARCH"
[ -z "$ARCH" ] && { # compatibility with openwrt versions prior to 19.07
    ARCH="$LEDE_ARCH"
}


[ -z "$ARCH" ] && {
  echo "Release information not available, skipping opkg configuration"
  exit 0
}

[ "$LIME_CODENAME" == "development" ] && {
	base_url="http://feed.libremesh.org/master";
	key_name="a71b3c8285abd28b"
	key_content="RWSnGzyChavSiyQ+vLk3x7F0NqcLa4kKyXCdriThMhO78ldHgxGljM/8"
} || {
	base_url="http://feed.libremesh.org/$LIME_RELEASE"
	key_name="a71b3c8285abd28b"
	key_content="RWSnGzyChavSiyQ+vLk3x7F0NqcLa4kKyXCdriThMhO78ldHgxGljM/8"
}

echo "Configuring official LibreMesh opkg feeds"
echo "src/gz libremesh $base_url/libremesh" > /etc/opkg/limefeeds.conf
echo "untrusted comment: signed by libremesh.org key $key_name" > "/etc/opkg/keys/$key_name"
echo "$key_content" >> "/etc/opkg/keys/$key_name"
