# packages/lime-system/files/etc/uci-defaults/95-lime-init-enable
#!/bin/sh

/etc/init.d/lime-init enable

