# packages/lime-system/files/etc/uci-defaults/99-lime-config
#!/bin/sh

sleep 5s
cp -r /etc/config /etc/lime/config.original
lua /usr/bin/lime-config && uci commit
/etc/init.d/lime-rc enable
