# packages/shared-state-dnsmasq_hosts/files/etc/uci-defaults/shared-state-dnsmasq_hosts-cron
#!/bin/sh

unique_append()
{
	grep -qF "$1" "$2" || echo "$1" >> "$2"
}

unique_append \
	'*/5 * * * * ((sleep $(($(hexdump -n 2 -e \"%u\" /dev/urandom) % 120)); shared-state sync dnsmasq-hosts &> /dev/null)&)'\
	/etc/crontabs/root
