# packages/shared-state-dnsmasq_hosts/files/etc/uci-defaults/shared-state-dnsmasq_hosts-cron
#!/bin/sh

. /lib/libremesh.sh

unique_append \
	'*/5 * * * * ((sleep $(($RANDOM % 120)); shared-state sync dnsmasq-hosts &> /dev/null)&)'\
	/etc/crontabs/root
