# packages/hotplug-initd-services/files/etc/uci-defaults/hotplug-initd-services
#!/bin/sh

mkdir -p /etc/hotplug.d/initd
