# packages/miniserver-client/files/etc/uci-defaults/90_miniserver-client
#!/bin/sh

uci set prometheus-node-push-influx.miniserver=prometheus-node-push-influx
uci set prometheus-node-push-influx.miniserver.server_address='miniserver'
uci set prometheus-node-push-influx.miniserver.server_port='8428'
uci set prometheus-node-push-influx.miniserver.interval='30'
uci set prometheus-node-push-influx.miniserver.disabled=
uci commit prometheus-node-push-influx
