# packages/deferrable-reboot/files/etc/uci-defaults/90-cron-defer-reboot
#!/bin/sh

unique_append()
{
	grep -qF "$1" "$2" || echo "$1" >> "$2"
}

unique_append \
	"*/3 * * * * ((sleep $((RANDOM % 120)); /usr/bin/deferrable-reboot-check )&)" /etc/crontabs/root

