# packages/shared-state/files/etc/uci-defaults/shared-state-publishers-cron
#!/bin/sh

unique_append()
{
	grep -qF "$1" "$2" || echo "$1" >> "$2"
}

unique_append \
	'*/30 * * * * ((sleep $(($RANDOM % 1000)); for publisher in /etc/shared-state/publishers/* ; do [ -x "$publisher" ] && "$publisher"; done )&)'\
	/etc/crontabs/root
