# packages/batman-adv-auto-gw-mode/files/etc/hotplug.d/net/99-batman-gw
#!/bin/sh

dhcp4_discover () {
	hostname=$(uci get system.@system[0].hostname)
	ubus call network.interface.lan_dhcp add_device '{ "name": "br-lan" }'
	cd /lib/netifd/proto/ && ./dhcp.sh dhcp setup lan_dhcp '{ "hostname": "'${hostname}'" }' br-lan
}

dhcp4_kill () {
	ubus call network.interface.lan_dhcp down
	ubus call network.interface.lan_dhcp remove_device '{ "name": "br-lan" }'
}

dnsmasq_start () {
	uci revert -P/var/state dhcp.@dnsmasq[0].domainneeded
	uci revert -P/var/state dhcp.lan.ignore
	/etc/init.d/dnsmasq start
}

dnsmasq_stop () {
	uci set -P/var/state dhcp.@dnsmasq[0].domainneeded=
	uci set -P/var/state dhcp.lan.ignore=1
	/etc/init.d/dnsmasq stop
}

if [ "$BATTYPE" = "gw" ] ; then
	case "$BATACTION" in
		add) dnsmasq_stop ; dhcp4_discover ;;
		del) dhcp4_kill ; dnsmasq_start ;;
		change) dhcp4_kill ; sleep 5 ; dhcp4_discover ;;
	esac
fi
