# packages/dnsmasq-lease-share/files/etc/uci-defaults/90_dnsmasq-lease-share
#!/bin/sh

uci set dhcp.@dnsmasq[0].dhcpscript=/usr/bin/dnsmasq-lease-share
#uci set dhcp.@dnsmasq[0].leasefilero=1
## Workaround until proper uci support is patched
echo leasefile-ro >> /etc/dnsmasq.conf
uci commit dhcp

exit 0
