# packages/dnsmasq-lease-share/files/etc/uci-defaults/90_dnsmasq-lease-share
#!/bin/sh

uci set dhcp.@dnsmasq[0].dhcpscript=/etc/alfred/dnsmasq-lease-share.lua
uci set dhcp.@dnsmasq[0].dhcphostsfile=/tmp/dhcp.hosts_remote
uci commit dhcp
uci set alfred.alfred.disabled=0
uci commit alfred

uci set alfred.alfred.disabled=
uci commit alfred

exit 0
