# packages/lime-system/files/etc/profile.d/10-lime-banner.sh
#!/bin/sh

echo "  ___   __ __                _______             __"
echo " |   |_|__|  |--.----.-----.|   |   |-----.-----|  |--."
echo " |     |  |  _  |   _|  -__||       |  -__|__ --|     |"
echo " |_____|__|_____|__| |_____||__|_|__|_____|_____|__|__|"
echo ""
echo " ------------------------------------------------------"
echo " %LIME_DESCRIPTION%"
echo " ------------------------------------------------------"
echo " https://libremesh.org"
echo " ------------------------------------------------------"
