# packages/lime-docs/files/etc/profile.d/15-lime-docs-banner
#!/bin/sh

echo " Find offline documentation in /docs or http://thisnode.info/docs/"
echo " ------------------------------------------------------"

