# packages/deferable-reboot/files/etc/uci-defaults/90-cron-defer-reboot
#!/bin/sh

unique_append()
{
	grep -qF "$1" "$2" || echo "$1" >> "$2"
}

unique_append \
	"*/3 * * * * ((if ping -c 10 -W 10 \$(uci -q get lime-autogen.system.deferable_reboot_ping_target) &> /dev/null; then awk '{print \$1 + 1800}' /proc/uptime > /tmp/deferable-reboot.defer ; fi )&) "\
	/etc/crontabs/root

