# packages/deferable-reboot/files/etc/uci-defaults/deferable-reboot-cron
#!/bin/sh

unique_append()
{
	grep -qF "$1" "$2" || echo "$1" >> "$2"
}

unique_append \
	'*/1 * * * * ((deferable-reboot &> /dev/null)&)'\
	/etc/crontabs/root
