# packages/deferable-reboot/files/etc/uci-defaults/deferable-reboot-cron
#!/bin/sh

. /lib/libremesh.sh

unique_append \
	'*/1 * * * * ((deferable-reboot &> /dev/null)&)'\
	/etc/crontabs/root
