# packages/eupgrade/files/etc/uci-defaults/99-eupgrades-cron
#!/bin/sh

unique_append()
{
	[ -f "$2" ] && grep -qF "$1" "$2" || echo "$1" >> "$2"
}

unique_append \
	'0 */6 * * * ((sleep $(($RANDOM % 120)); eupgrade-check &> /dev/null)&)'\
	/etc/crontabs/root
