# packages/lime-smart-wifi/files/etc/uci-defaults/90_lime-smart-wifi
#!/bin/sh

uci -q get lime-community.smart_wifi || \
    uci set lime-community.smart_wifi='smart_wifi'
uci -q get lime-community.smart_wifi.ap_2ghz || \
    uci set lime-community.smart_wifi.ap_2ghz=1
uci -q get lime-community.smart_wifi.ap_5ghz || \
    uci set lime-community.smart_wifi.ap_5ghz=0
uci -q get lime-community.smart_wifi.mesh_2ghz || \
    uci set lime-community.smart_wifi.mesh_2ghz=1
uci -q get lime-community.smart_wifi.mesh_5ghz || \
    uci set lime-community.smart_wifi.mesh_5ghz=1
uci -q get lime-community.smart_wifi.all_to_mesh || \
    uci set lime-community.smart_wifi.all_to_mesh=0
uci -q get lime-community.smart_wifi.all_to_ap || \
    uci set lime-community.smart_wifi.all_to_ap=0
uci -q get lime-community.smart_wifi.channels_2ghz || \
    uci add_list lime-community.smart_wifi.channels_2ghz="1 11 6"
uci -q get lime-community.smart_wifi.channels_5ghz || \
    uci add_list lime-community.smart_wifi.channels_5ghz="36 48"
uci commit lime-community

/usr/bin/lime-smart-wifi > /tmp/log/lime-smart_wifi.log

uci commit lime
