# packages/ubus-lime-utils/files/etc/uci-defaults/90-create_udhcpc.user
#!/bin/sh

### this is needed for compatibility with old OpenWrt 19.07

[ -f /etc/udhcpc.user ] || cat <<EOF > /etc/udhcpc.user
#!/bin/sh
for f in /etc/udhcpc.user.d/*; do
	[ -f "\$f" ] && (. "\$f" "\$@")
done

EOF
