# packages/shared-state-network_nodes/files/etc/uci-defaults/shared-state-network_nodes-cron
#!/bin/sh

. /lib/libremesh.sh

unique_append \
	'*/5 * * * * ((sleep $(($RANDOM % 120)); shared-state-multiwriter sync network_nodes &> /dev/null)&)'\
	/etc/crontabs/root
