# packages/shared-state-dnsmasq_leases/files/etc/uci-defaults/90_dnsmasq-lease-share
#!/bin/sh

dhcpHostsFile="/tmp/dhcp.hosts_remote"

uci set dhcp.@dnsmasq[0].dhcpscript=/usr/bin/dnsmasq-lease-share.sh
uci set dhcp.@dnsmasq[0].dhcphostsfile=$dhcpHostsFile
uci commit dhcp

unique_append()
{
	grep -qF "$1" "$2" || echo "$1" >> "$2"
}

unique_append \
	'*/5 * * * * ((sleep $(($(hexdump -n 2 -e \"%u\" /dev/urandom) % 120)); shared-state sync dnsmasq-leases &> /dev/null)&)'\
	/etc/crontabs/root

exit 0
