# packages/shared-state-dnsmasq_leases/files/etc/uci-defaults/90_dnsmasq-lease-share
#!/bin/sh

uci set dhcp.@dnsmasq[0].dhcpscript=/usr/bin/dnsmasq-lease-share.sh
uci set dhcp.@dnsmasq[0].dhcphostsfile=/tmp/dhcp.hosts_remote
uci commit dhcp

exit 0
