# packages/check-date-http/files/etc/uci-defaults/check-date-http-cron
#!/bin/sh

unique_append()
{
	grep -qF "$1" "$2" || echo "$1" >> "$2"
}

unique_append \
	'*/20 * * * * ((sleep $(($RANDOM % 600)); check-date-http &> /dev/null)&)'\
	/etc/crontabs/root
