# packages/check-date-http/files/etc/uci-defaults/check-date-http-cron
#!/bin/sh

. /lib/libremesh.sh

unique_append \
	'*/20 * * * * ((sleep $(($RANDOM % 600)); check-date-http &> /dev/null)&)'\
	/etc/crontabs/root
