# packages/lime-system/files/etc/init.d/lime-init
#!/bin/sh /etc/rc.common

START=89

start()
{
	local hooks=/etc/lime-init.d/*.start
	for hook in $hooks ; do
	. $hook || true
	done
}

stop()
{
	local hooks=/etc/lime-init.d/*.stop
	for hook in $hooks ; do
	. $hook || true
	done
}

restart()
{
	stop
	sleep 3
	start
}
