# packages/lime-ap-watchping/files/etc/uci-defaults/91_autoAP
#!/bin/sh
uci -q get lime-autogen.network.autoap_enabled || uci set lime-autogen.network.autoap_enabled=0
uci -q get lime-autogen.network.autoap_hosts || uci set lime-autogen.network.autoap_hosts="8.8.8.8 8.8.4.4"
