# packages/first-boot-wizard/files/etc/uci-defaults/zzzz-serve-lime-config
#!/bin/sh

ln -s /etc/config/lime-defaults /www/lime-defaults
