# packages/smonit/files/etc/init.d/smonit
#!/bin/sh /etc/rc.common

START=99
SMONIT="/usr/sbin/smonit"
CRONTAB="/etc/crontabs/root"
CRON_LINE='*/5 * * * * $SMONIT'

start() {
  echo "Ensuring smonit is run by cron"
  if ! ( grep -q "$CRON_LINE" "$CRONTAB" 2>/dev/null ) ; then
    echo "$CRON_LINE" >> "$CRONTAB"
    /etc/init.d/cron enable
    /etc/init.d/cron restart
  fi
}

stop() {
  echo "Removing smonit from crontab"
  sed "\|$CRON_LINE|d" -i "$CRONTAB"
  /etc/init.d/cron restart
}

reload() {
  $SMONIT
}
