# packages/smonit/files/etc/init.d/smonit
#!/bin/sh /etc/rc.common

START=99
SMONIT="/usr/sbin/smonit"
CRONTAB="/etc/crontabs/root"
CRON_LINE='*/5 * * * *'

start() {
  echo "Enabling smonit to crontab"
  grep -q "$SMONIT" $CRONTAB || echo "$CRON_LINE $SMONIT" >> $CRONTAB
  /etc/init.d/cron reload
}

stop() {
  echo "Disabling smonit from crontab"
  cat $CRONTAB | grep -v "$SMONIT" > /tmp/cron.tmp
  cp -f /tmp/cron.tmp $CRONTAB
}

reload() {
  $SMONIT
}
