# packages/shared-state-persist/files/etc/uci-defaults/50-set-automount
#!/bin/sh
uci set fstab.@global[0].auto_mount='1'
uci set fstab.@global[0].anon_mount='1'
uci commit fstab
