# packages/shared-state-babeld_hosts/files/etc/hotplug.d/initd/shared-state-babeld_hosts
#!/bin/sh
# SPDX-License-Identifier: AGPL-3.0-only
#
# Copyright (C) 2019  Gioacchino Mazzurco <gio@altermundi.net>

[ "$SERVICE" == "babeld" ] && [ "$ACTION" == "service.start" ] &&
	((sleep 2s; /usr/bin/shared-state-publish_babeld_hosts)&)
