# packages/lime-firewall/files/etc/uci-defaults/95-lime-firewall-enable
#!/bin/sh

/etc/init.d/lime-firewall enable

