# packages/lime-docs/files/etc/uci-defaults/91_lime-docs-banner
#!/bin/sh

cat <<EOF >> /etc/banner
 Find offline documentation in /docs or http://thisnode.info/docs/
 ------------------------------------------------------
EOF

exit 0
