# packages/shared-state/files/etc/init.d/shared-state
#!/bin/sh /etc/rc.common

START=98

start() {
    logger -p daemon.info  -t shared-state "Generating and publishing databases"
    shared-state-publish-all &> /dev/null
}
