# packages/tests-scheduler/files/etc/init.d/tests-scheduler-probe-data-persist
#!/bin/sh /etc/rc.common
# Copyright (C) 2019 Ilario Gelmetti <ilario@sindominio.net>
USE_PROCD=1
START=45
STOP=01
CONFIGURATION=tests-scheduler

start_service() {
	config_load "${CONFIGURATION}"
	local persist_to_flash
	local data_file
	local backup_dir
	local backup_filename
	config_get persist_to_flash probe persist_to_flash
	config_get data_file probe data_file
	config_get backup_dir probe backup_dir
	config_get backup_filename probe backup_filename
	
	procd_open_instance
	procd_set_param respawn
	procd_set_param command $persist_to_flash && [[ -e $backup_dir/$backup_filename ]] && cp $backup_dir/$backup_filename $data_file
	procd_set_param file /etc/config/tests-scheduler
	procd_set_param stdout 1
	procd_set_param stderr 1
	procd_close_instance
}

stop_service() {
	config_load "${CONFIGURATION}"
	local persist_to_flash
	local data_file
	local backup_dir
	local backup_filename
	config_get persist_to_flash probe persist_to_flash
	config_get data_file probe data_file
	config_get backup_dir probe backup_dir
	config_get backup_filename probe backup_filename
	
	$persist_to_flash && [[ -e $data_file ]] && [[ -d $backup_dir ]] && cp $data_file $backup_dir/$backup_filename
}
