# packages/soporteremoto/files/etc/uci-defaults/95_soporteremoto_sshkeys
#!/bin/sh
sshkeys="
ssh-rsa AAAAB3NzaC1yc2EAAAADAQABAAACAQDBK9+aREyRQT4FRXTtYP+5Uke0rzzcFosDAQ6nzCitNWrZcMDWuL89mEvuP66HeGkI8HivgbZKEMgxdZgBkElFcoEmr0YqbyBI+eI580dCNHhSXkT6xFu+0kZV2vrvEF2/eocAuhwTO80NQW3vx3HS/xCJYbwriE8f9KxK/aFQhLaUD8p9ANUHDJUtEFAD/Ma3rAvS+peeKUUYkzqqI29gK4cFvUo7YF+6M7A19jgyk6B5q5pGU4kPeopdx1a4sIGa4w1+vvTCTSNKz9K6xiF4Wn1XW98nD83uXFPfPII4EuBAGw7D4KpwHr2EhIOHHic426YAHQwmTWhvk8sZfh2mCJ706weg3jLd+Skk71S9RlLfwtw+ExsmwqR6iY00NjGm1ccSpL95hEISra4k14WhxXqUPd36zE0BiYi868vRnWxl5yXvC6la/EqwgJHI12m2LRHzR0WKbiLcB1sCAIKUKEifweV0XHkEcuzyIWlLMRVMH+zfAHMeWqxMrQjVBmEmI+xk9uK+JqZqKf34fcdCZUo247Updg9ass1KJgEaVozGRwezXju+Dwv3fs+ZHWKT2PSDgh04GnimbukQ/PmU6ht6tYIB000KAUYcec0/TNcLR6Vnr4tqynNmlSmJi0BbUI0lauzCr5zCiqZ9Lkx+EGzcl1nR7Fbstir8qN2oeQ== soporteremoto@altermundi.net
ssh-rsa AAAAB3NzaC1yc2EAAAABIwAAAQEAp9Cb57jGm0MaN1x+1VmS7BFDo7+BCKP+5Kl0JSvzCwTbrYlUUDZLQ/GHkm6e/4irJWIxqGTfEOOc/VH1AMiuVIN23tp81itEuoZWqcqJRb81CGVi8zyuTJ/F1+WVOjUtYzjzaY2ZllVphZIWLOCjAakcEZzrqh+lESDfslnX17ezvCcDduBUrNmaZ0F2rJipbip7pda62v4SOQsV1MoPmOwcQeLqyLbDzQOVFRKAqZUZZ4hk+wlBtGU5PwDvU5SGJGbR9tZmsE5P+fNYUMnr4Ri0DHbKPCmygOT2J6zJ9hcil3Vct/t/mOTJWryTarjbVz74f/Uf4Ig8q3NegbGEYw== id_guido@laptop
ssh-rsa AAAAB3NzaC1yc2EAAAABIwAAAQEAp/zdIgomuLuNevi1wIYVPNiGjYrPdx8Fpw2H3b87TCRD90F9NEpyZWYWAEYlPioErR+RBlJfq+TEyxvagRDxh+u5GDmggRAF7tdOPe2QNLWdcrHNoUQHghoZKP70Qp1rSc6u2/YUs0CiHJ+VjeM6iR+HJuheJQ/myF36OqNVLsx+uHUwIW6QM/9zxVdyQMo3Dgb9HJ/wMzc1YbqwF19ANZ+n/GdzmNc6/6ng4SxZAVet0xBcdKe5/9wlpUpjehxiAJNKrtQUgOE6YBAFxy69yKxluupZRBpMSWYAzC1jytXfNflqXSt4Eq448OBYnRbxOwVuO/rgohNmAMlrcw+psw== nicolas@fierrito
ssh-rsa AAAAB3NzaC1yc2EAAAADAQABAAACAQC8jbnSTO+i83zlwoSGYEvj0dFrNVbuiEW7Y+9GnbSsJ6Guh5acIeu3ySOHVqXB4+hrcPF/pKqRf0bkz2D+ccBmtFhjxFxOwlAYpBJpU+Lwgvw5+2K8s4f0Bq+grY8GMrsByk5g9R+3ykVeLqhAMLOrFRPo5P0Pb+Nwvpjk8ZRbR6S3kblAvhk7n9ejT3G8VgKlLi93Kg3ZGzuo96sLchEKXeLI208EApZ+CrqkMycAPcrpHpIee70JgL0gDuqoxbUhMkdlQNPi7Fz6XGJCp4pFmOoCXA5P+n/VhXY3x4gcIcRZoTRgzhsbHZtIVpbOZ2zhBShN3pKHL0tm/M1OiVJEi/gMG8Wh0iDlg6xdC4nJcxyliYJha6Ca1fJCsh9FAQu51PuJwPzDDsruV78LzvJxtzp39DUw7u2lYsVYmUod9xSNhL45qa88qWq2WqEzBmkLbatef/4tsjoOd3KrqoFNdQaMf8Yc1fi6zEgiID+IvZPW38iJwVKB1oLHV+7k8bINrNhXLrtrSsM/YS1x2RKP/jDVRb+jQ7ZGpsKSmJaos3OMrUunJLsJK7m0DnSnN+CC0BhbKhdyUUcojwvQqWmLaLtIGIQ3Ovx/ovmtvKVK/FqaJdrj7rOM7ETKSvZBz4DF8+D7r3eJtx4ks7u41bExJM6su1Iw5nQgo3oOh7fiOw== Gio <gio@diveni.re>
ssh-rsa AAAAB3NzaC1yc2EAAAADAQABAAABAQDDmmm+KDqSsZkvgbY3Mm2Zdr+0KmVWjBPdMMktc1tmvna28xP/kUXfzy4Z0JStkCBl/9wWXe7QYPZpoyoL2i9uSrFlZoWRxvcGS/HGCnS7F33YCyoSNnojR8nSiW2KDmLLODIEUsh/btbgQXbo92hDMAgJIPl/GckQeRn2RqxUoOdOhXD1OJ0q9Zpc3wSHQfPd9i7hYztYs5+J/pFiBU/q/71B2VWCdjC8EwzdJmEEErx6h3eX7MMFITNtQtXjsgCkFUDlLQ6xI3WMmEF9ks6Bgw60C0LCsjXRC1EZtba2r26agrPKKWfjqMjcLHwgUIaot5oFjexgvv17JuyJsV6J san@bonzo
"

authkeys="/etc/dropbear/authorized_keys"
umask 077 ; mkdir -p $(dirname $authkeys) ; touch $authkeys
IFS="
"
for key in $sshkeys ; do
  sed "\|^$key$|d" -i $authkeys
  echo   "$key"    >> $authkeys
done

### Disable root SSH access via password
passwd -l root
