# packages/cotonete/files/etc/init.d/cotonete
#!/bin/sh /etc/rc.common

START=99
USE_PROCD=1
PROG=/usr/sbin/cotonete.sh

start_service() {
  procd_open_instance [cotonete]
  procd_set_param command $PROG

  # respawn automatically if something died, be careful if you have an alternative process supervisor
  # if process dies sooner than respawn_threshold, it is considered crashed and after 5 retries the service is stopped
  procd_set_param respawn ${respawn_threshold:-3600} ${respawn_timeout:-5} ${respawn_retry:-5}

  procd_set_param stdout 1 # forward stdout of the command to logd
  procd_set_param stderr 1 # same for stderr
  procd_set_param user root
  procd_set_param pidfile /var/run/cotonete.pid # write a pid file on instance start and remote it on stop
  procd_close_instance
}
